# packaging/efl-config.path
[Unit]
Description=EFL Configuration Daemon

[Path]
PathExists=/run/.wm_ready
