# packaging/efl-config.service.in
[Unit]
Description=Change the efl configures.
After=enlightenment-user.service

[Service]
Type=simple
ExecStart=@TZ_SYS_BIN@/efl_config
EnvironmentFile=-/run/tizen-system-env

[Install]
WantedBy=default.target
